# misc/rkg.cron
# This is cron entry every 5 minutes for rkg.php. You have to adjust your path.
*/5 * * * * www-data cd /var/www/wpa-sec && /usr/bin/flock -n /tmp/rkg.php.lock /usr/bin/php rkg.php

