# misc/rkg.cron
# This is cron entry every 5 minutes for rkg.php. You have to adjust your path.
* * * * * cd /var/www/html && /usr/bin/flock -n /tmp/rkg.php.lock /usr/local/bin/php rkg.php
