# misc/services.conf
[supervisord]
nodaemon=true
user=root

[program:apache]
command=/usr/local/bin/apache2-foreground

[program:cron]
command=/usr/sbin/cron -f